# eap_cluster/docker/eap/node1/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:docker_init]
command=/docker_init.sh
